# grafana/roles/grafana/files/etc/systemd/system/grafana-server.service.d/override.conf
# NOTE(mastensg): for binding to port < 1024.
[Service]
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
PrivateUsers=false
