# shared-roles/pdb/templates/crontab
45 03 * * * root /opt/wal-g/wal-g backup-push -f $PGDATA
